# docker/serv/configs/crontab
* * * * * /bin/python3 /usr/src/event-monitor-snort3/snort3_monitor/clear_deleted_events.py

* * * * * /bin/python3 /usr/src/event-monitor-snort3/snort3_monitor/mongo/perf_monitor_listener.py
